# 02_recursos/censoestudiantil.socket
[Unit]
Description=Socket de la API de Censo Estudiantil

[Socket]
ListenStream=/run/censoestudiantil.sock

[Install]
WantedBy=sockets.target